# lib/systemd/system/mega.service
[Unit]
Description=Mega Provider Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mega


ExecStart=/opt/mega/venv/bin/python /opt/mega/server.py --port=6996

ExecStartPost=/bin/bash -c "nohup /usr/lib/jvm/java-11/bin/java -cp /opt/mega/ComputeKart-tunnel-client-1.0.0.jar main.Main ${TUNNEL_SERVER_TOKEN} localhost 6996 > /opt/mega/tunnel.log 2>&1 &"

Restart=always
EnvironmentFile=/etc/mega/config.conf

[Install]
WantedBy=multi-user.target
